# gridcoinresearchd/gridcoinresearchd.service
[Unit]
Description=Gridcoin daemon

[Service]
Type=forking
ExecStart=/usr/bin/gridcoinresearchd
User=edwinksl

[Install]
WantedBy=multi-user.target
